# initramfs/etc/inittab


::sysinit:/etc/init.d/rcS

::askfirst:-/bin/login

tty2::askfirst:/bin/login
tty3::askfirst:/bin/login
tty4::askfirst:/bin/login

tty4::respawn:/sbin/getty 38400 tty5
tty5::respawn:/sbin/getty 38400 tty6


::restart:/sbin/init

::ctrlaltdel:/sbin/reboot

::shutdown:/bin/umount -a -r
::shutdown:/sbin/swapoff -a

